# server/ipcam.service
# ipcam.service — systemd 服务单元文件
#
# 安装步骤：
#   sudo cp ipcam.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable ipcam
#   sudo systemctl start ipcam
#
# 查看日志：
#   sudo journalctl -u ipcam -f

[Unit]
Description=Pisces IP Camera Server
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=pi
Group=pi
WorkingDirectory=/opt/ipcam-server

# 从 .env 文件加载环境变量（生产环境也可直接在此设置 Environment=）
EnvironmentFile=/opt/ipcam-server/.env

ExecStart=/opt/ipcam-server/venv/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --log-level info \
    --no-access-log

# 异常退出后 5 秒自动重启
Restart=on-failure
RestartSec=5s

# 资源限制
LimitNOFILE=65536

# 标准输出/错误均交给 journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipcam

[Install]
WantedBy=multi-user.target
